# deploy/systemd/domeneshop-mcp.service.example
[Unit]
Description=Domeneshop MCP Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/opt/domeneshop-mcp
ExecStart=/opt/domeneshop-mcp/.venv/bin/domeneshop-mcp-server
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
